# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:create-branch]
command=/golang/pinger/bin/pinger.sh
