# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:create-branch]
command=source /pinger-conf/env.sh;/golang/bin/pinger -config=/pinger-conf/pinger.yml